# supervisord.conf
[supervisord]
nodaemon=true

[program:laravel-scheduler]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan schedule:run --verbose --no-interaction
autostart=true
autorestart=true
startsecs=0
stdout_logfile=/var/log/supervisor/laravel-scheduler.log
stderr_logfile=/var/log/supervisor/laravel-scheduler-error.log
user=www-data
numprocs=1

[program:laravel-server]
process_name=%(program_name)s
command=php /var/www/artisan serve --host=0.0.0.0 --port=8000
autostart=true
autorestart=true
startsecs=0
stdout_logfile=/var/log/supervisor/laravel-server.log
stderr_logfile=/var/log/supervisor/laravel-server-error.log
user=www-data
